# life-systems.service
[Unit]
Description=Life Systems FastAPI Application
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/opt/life-systems
Environment="PATH=/opt/life-systems/venv/bin"
Environment="LS_USER=jurek"
Environment="LS_PASSWORD=LifeSystems2026!"
EnvironmentFile=-/opt/life-systems/.env
ExecStart=/opt/life-systems/venv/bin/uvicorn api.main:app --host 127.0.0.1 --port 8000
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
