# systemd/evening-digest.timer
[Unit]
Description=Life Systems - Evening Activities Digest Timer
Documentation=https://github.com/kevin-bot-openclaw-ops/life-systems

[Timer]
# Run at 20:00 UTC (21:00 CET in winter, 22:00 CEST in summer)
OnCalendar=*-*-* 20:00:00
Persistent=true

[Install]
WantedBy=timers.target
