# systemd/life-systems-scanner.service
[Unit]
Description=Life Systems Job Scanner
After=network.target

[Service]
Type=oneshot
User=ubuntu
WorkingDirectory=/opt/life-systems
Environment="PATH=/opt/life-systems/venv/bin:/usr/local/bin:/usr/bin:/bin"
Environment="DB_PATH=/var/lib/life-systems/life.db"
Environment="PYTHONPATH=/opt/life-systems"
EnvironmentFile=/etc/life-systems/env
ExecStart=/opt/life-systems/venv/bin/python3 -c "import asyncio; from api.database import Database, init_db; from api.scanner import run_scan; db_path='/var/lib/life-systems/life.db'; init_db(db_path); db=Database(db_path); asyncio.run(run_scan(db))"
