# systemd/life-systems.service
[Unit]
Description=Life Systems API
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/opt/life-systems
Environment="PATH=/opt/life-systems/venv/bin:/usr/local/bin:/usr/bin:/bin"
Environment="DB_PATH=/var/lib/life-systems/life.db"
Environment="PORT=8000"
Environment="PYTHONPATH=/opt/life-systems"
EnvironmentFile=/etc/life-systems/env
ExecStart=/opt/life-systems/venv/bin/uvicorn api.main:app --host 127.0.0.1 --port 8000
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
